# rootfs/etc/crontabs/caddy
0 8 * * * /usr/bin/php -f /srv/www/artisan auth:clear-reminders > /dev/null 2>&1
