# latest/overlay/etc/crontabs/caddy
0 8 * * * cd /srv/www && /usr/bin/php -f artisan auth:clear-reminders > /dev/null 2>&1
